# z_archive/a_kinoite_host/1_oci/scripts/configure-sessions.service
[Unit]
Description=Configure desktop sessions based on profile
After=local-fs.target
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/configure-sessions.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
